# templates/ubuntu/etc/profile.d/php-aliases.sh.erb
alias php-debug='XDEBUG_CONFIG="profiler_enable=0" PHP_IDE_CONFIG="serverName=${HOSTNAME}.local" php'
alias php-profile='XDEBUG_CONFIG="profiler_enable=1" PHP_IDE_CONFIG="serverName=${HOSTNAME}.local" php'
